# storage-appliance/lib/systemd/system/petasan-iscsi.service
[Unit]
Description= PetaSAN iSCSI Service

[Service]
Type=simple
ExecStart=/opt/petasan/services/iscsi_service.py
Restart=always
RestartSec=3
