# storage-appliance/lib/systemd/system/petasan-admin.service
[Unit]
Description= PetaSAN Web Management and Administration

[Service]
Type=simple
ExecStart=/opt/petasan/services/web/admin.py
Restart=always
RestartSec=3
