# storage-appliance/lib/systemd/system/petasan-file-sync.service
[Unit]
Description= PetaSAN File Sync Service

[Service]
Type=simple
ExecStart=/opt/petasan/services/files_sync.py
Restart=always
RestartSec=3
